# xyz.openbmc_project.scratchpad.service
[Unit]
Description=scratchpad
After=dbus.service

[Service]
ExecStart=/usr/sbin/scratchpad
Type=simple

[Install]
WantedBy=basic.target
